# debian/pop-upgrade.postinst
#!/bin/sh
set -x
systemctl --global enable nux-upgrade-notify.timer
systemctl disable nux-upgrade

if ! test -e /var/lib/nux-upgrade/restarting && ! test -e /nux-upgrade && ! test -e /pop_preparing_release_upgrade; then
    if systemctl is-active nux-upgrade; then
        systemctl restart nux-upgrade
    fi
fi

#DEBHELPER#

exit 0
